# deploy/systemd/weft-network.service
# weft-network.service — systemd unit for the network control-plane daemon.
#
# Install :
#   - copy the weft-network binary to /usr/local/bin/weft-network
#   - copy this file to /etc/systemd/system/weft-network.service
#   - create the runtime + state dirs : `weft-network-prepare` (below)
#   - systemctl daemon-reload && systemctl enable --now weft-network
#
# Configure via /etc/default/weft-network — see the EnvironmentFile
# directive below. Key knobs :
#   WEFT_NETWORK_LISTEN — defaults to a unix socket under /run.
#   WEFT_NETWORK_ETCD   — empty = in-memory dev ; set for HA.
#   WEFT_NETWORK_LOG    — debug / info / warn / error.

[Unit]
Description=openweft network control-plane (routers / LBs / DNS / scheduling)
Documentation=https://openweft.github.io/docs/latest/
After=network-online.target etcd.service
Wants=network-online.target

[Service]
Type=simple
User=weft
Group=weft

# Defaults can be tuned in /etc/default/weft-network without editing
# this unit file. Missing file = use the daemon's own defaults.
EnvironmentFile=-/etc/default/weft-network

Environment="WEFT_NETWORK_LISTEN=unix:///run/weft-network/weft-network.sock"
# Bind to localhost ; expose via WireGuard / loadbalancer in production
# so Prometheus scrapes don't leak outside the cluster mesh.
Environment="WEFT_NETWORK_METRICS=127.0.0.1:9100"
Environment="WEFT_NETWORK_ETCD="
Environment="WEFT_NETWORK_LOG=info"
# TLS — only meaningful for TCP listeners. Leave empty when listening
# on the default unix socket (filesystem perms gate access). Set all
# three to require mTLS ; cert+key without client-ca = TLS, anonymous
# clients.
Environment="WEFT_NETWORK_TLS_CERT="
Environment="WEFT_NETWORK_TLS_KEY="
Environment="WEFT_NETWORK_CLIENT_CA="
# Tracing — OTLP/gRPC collector endpoint. Empty = disabled (default).
# Point at Tempo, Jaeger's OTLP receiver, or any OTLP-compatible vendor.
# Example : WEFT_NETWORK_OTLP_ENDPOINT=tempo.weft.internal:4317
Environment="WEFT_NETWORK_OTLP_ENDPOINT="

# The pre-start step prepares /run/weft-network with the right perms.
# RuntimeDirectory= would also work but ExecStartPre is portable to
# older systemd that doesn't ship the same RuntimeDirectoryMode logic.
RuntimeDirectory=weft-network
RuntimeDirectoryMode=0750
StateDirectory=weft-network
StateDirectoryMode=0750

ExecStart=/usr/local/bin/weft-network \
  --listen        "${WEFT_NETWORK_LISTEN}" \
  --metrics-addr  "${WEFT_NETWORK_METRICS}" \
  --etcd          "${WEFT_NETWORK_ETCD}" \
  --log-level     "${WEFT_NETWORK_LOG}" \
  --tls-cert      "${WEFT_NETWORK_TLS_CERT}" \
  --tls-key       "${WEFT_NETWORK_TLS_KEY}" \
  --client-ca     "${WEFT_NETWORK_CLIENT_CA}" \
  --otlp-endpoint "${WEFT_NETWORK_OTLP_ENDPOINT}"

# Cooperative reload : SIGHUP is reserved for future config re-read.
# Today only SIGTERM is honoured (daemon GracefulStop's the gRPC server).
KillMode=mixed
KillSignal=SIGTERM
TimeoutStopSec=30s
Restart=on-failure
RestartSec=5s

# Hardening — tight enough for a stateless daemon. Loosen only if
# you wire a local-fs backend or attach to a custom path.
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
PrivateTmp=true
PrivateDevices=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
RestrictNamespaces=true
RestrictRealtime=true
RestrictSUIDSGID=true
LockPersonality=true
MemoryDenyWriteExecute=true
SystemCallArchitectures=native
SystemCallFilter=@system-service
SystemCallErrorNumber=EPERM

[Install]
WantedBy=multi-user.target
